# rules/persistence/credential_access_modify_ssh_binaries.yaml
# Title: Credential Access Modify Ssh Binaries

# Creation date: 2020/12/21

# MITRE ATT&CK Tactic: TA0003 - Persistence - https://attack.mitre.org/tactics/TA0003/

# MITRE ATT&CK Technique: T1543 - Create or Modify System Process - https://attack.mitre.org/techniques/T1543/

- name: Detect unauthorized SSH binary modification
  type: FileOpened
  category: persistence
  severity: medium
  description: Detects changes to SSH-related binaries by non-package management processes.
    Adversaries may modify SSH-related binaries for persistence or credential access by patching sensitive functions to enable unauthorized access or by logging SSH
    credentials for exfiltration.
  condition: (payload.filename IN [
      "/usr/bin/scp", 
      "/usr/bin/sftp", 
      "/usr/bin/ssh", 
      "/usr/sbin/sshd", 
      "/libkeyutils.so", 
      "/usr/lib/openssh/sftp-server",
      "/usr/libexec/openssh/sftp-server",
      "/usr/lib/systemd/system/sshd.service",
      "/lib/x86_64-linux-gnu/libkeyutils.so.1",
      "/usr/lib/libkeyutils.so.1"
    ])
    AND (payload.flags CONTAINS "O_WRONLY" OR payload.flags CONTAINS "O_RDWR")
    AND NOT header.image IN [
      "/usr/bin/dnf", 
      "/usr/bin/dnf-automatic", 
      "/usr/bin/dpkg", 
      "/usr/bin/yum", 
      "/usr/bin/rpm", 
      "/usr/bin/yum-cron", 
      "/usr/bin/anacron"
    ]
