# rules/persistence/credential_access_modify_ssh_binaries.yaml
# Title: Credential Access Modify Ssh Binaries

# Creation date: 2020/12/21

# MITRE ATT&CK Tactic: TA0003 - Persistence - https://attack.mitre.org/tactics/TA0003/

# MITRE ATT&CK Technique: T1543 - Create or Modify System Process - https://attack.mitre.org/techniques/T1543/

- name: Detect unauthorized SSH binary modification
  type: FileOpened
  category: persistence
  severity: high
  description: Detects changes to SSH related binaries by non-package management processes. Adversaries may modify SSH related 
    binaries for persistence or credential access by patching sensitive functions to enable unauthorized access or by logging SSH 
    credentials for exfiltration.
  condition: (payload.filename IN ["/usr/bin/scp", "/usr/bin/sftp", "/usr/bin/ssh", "/usr/sbin/sshd", "/libkeyutils.so"]) 
    AND NOT header.image IN ["/usr/bin/dnf", "/usr/bin/dnf-automatic", "/usr/bin/dpkg", "/usr/bin/yum", "/usr/bin/rpm", 
    "/usr/bin/yum-cron", "/usr/bin/anacron"] 
  

# MITRE ATT&CK Technique: T1543 - Create or Modify System Process: Systemd Service - https://attack.mitre.org/techniques/T1543/002

# PROBLEMS with .config, not triggering. there is a way to get user name
- name: Detect systemd services modification
  type: FileOpened
  category: persistence
  severity: medium
  description: Detect creation or modification of systemd service files in /etc or /lib by non-package manager processes. 
    Adversaries may create or modify systemd services to repeatedly execute malicious payloads as part of persistence.
  condition: (payload.filename ENDS_WITH ".service")
    AND (
      payload.filename STARTS_WITH "/etc/systemd/system" OR
      payload.filename STARTS_WITH "/etc/systemd/user" OR
      payload.filename STARTS_WITH "/lib/systemd/system" OR
      payload.filename STARTS_WITH "/usr/lib/systemd/system" OR
      payload.filename STARTS_WITH "/usr/lib/systemd/user" OR
      payload.filename STARTS_WITH "/usr/local/lib/systemd/system" OR
      payload.filename STARTS_WITH "/root/.config/systemd/user/" OR
      payload.filename STARTS_WITH "/root/.local/share/systemd/user" OR
      payload.filename STARTS_WITH ".config/systemd/" OR
      payload.filename STARTS_WITH "~/.config/systemd/"
    )
    AND NOT header.image IN [
      "/usr/bin/dnf", "/usr/bin/dnf-automatic", "/usr/bin/dpkg", "/usr/bin/yum", "/usr/bin/rpm", 
      "/usr/libexec/platform-python", "/usr/bin/gpg", "/usr/bin/yum-cron", "/usr/bin/anacron", 
      "/usr/lib/systemd/systemd", "/usr/lib/systemd/systemd-userwork"
    ]